# .crontab
@ reboot rm -rf "$HOME/.cache/" 2>/dev/null
@ reboot rm -rf "$HOME/.thumbnails" 2>/dev/null
@ reboot rm -rf "$HOME/.local/share/Trash" 2>/dev/null
@reboot rm -rf $HOME/.wine/drive_c/windows/temp/* 2>/dev/null
